# netbird-monitor-1.0/DEBIAN/postinst
#!/bin/bash
# netbird-monitor-1.0/DEBIAN/postinst

set -e

case "$1" in
    configure)
        # Script ausführbar machen
        chmod +x /usr/local/bin/netbird-monitor.sh

        # Log-Datei erstellen
        touch /var/log/netbird-monitor.log
        chmod 644 /var/log/netbird-monitor.log

        # Logrotate-Konfiguration aktivieren
        chmod 644 /etc/logrotate.d/netbird-monitor

        # Systemd-Services neu laden
        systemctl daemon-reload

        # Timer aktivieren und starten
        systemctl enable netbird-monitor.timer
        systemctl start netbird-monitor.timer

        echo "NetBird Monitor wurde installiert und gestartet."
        echo "Status prüfen mit: systemctl status netbird-monitor.timer"
        echo "Logs anzeigen mit: tail -f /var/log/netbird-monitor.log"
        echo "Logrotate-Test mit: logrotate -d /etc/logrotate.d/netbird-monitor"
        ;;
esac

exit 0
