# src/netbird-monitor.service
# /etc/systemd/system/netbird-monitor.service
[Unit]
Description=NetBird Connection Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/netbird-monitor.sh
Environment="MIN_PEERS=3"
Environment="STATUS_TIMEOUT=10"
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
